# backend/build/web_supervisor.conf
[program:web-worker]
directory=/opt/rnatango
command=gunicorn -b :8000 rnatango.asgi:application -w 4 -k uvicorn.workers.UvicornWorker -t 30
autostart=true
autorestart=true
user=root
numprocs=1
startretries=50
stopsignal=INT
killasgroup=true
stopasgroup=true
stdout_logfile=/opt/rnatango/supervisor/logs.conf
redirect_stderr=true
