# backend/build/ws_supervisor.conf
[program:ws-worker]
directory=/opt/rnatango
command=uvicorn rnatango.asgi:application --reload --host 0.0.0.0 --port 8081
autostart=true
autorestart=true
user=root
numprocs=1
startretries=50
stopsignal=INT
killasgroup=true
stopasgroup=true
stdout_logfile=/opt/rnatango/supervisor/logs.conf
redirect_stderr=true
